# deploy/systemd/novapilot-backend.service.example
[Unit]
Description=NovaPilot AI FastAPI Backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/novapilot/backend
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/opt/novapilot/backend/.env
ExecStart=/opt/novapilot/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
